# config/systemd/user/battery-monitor.service
[Unit]
Description=DWM Battery Monitor Check

[Service]
Type=oneshot
ExecStart=dwm-battery-monitor
LogLevelMax=warning
